# .config/uwsgi/uwsgi.service
[Unit]
Description=uWSGI service
After=syslog.target
 
[Service]
ExecStart=/home/ubuntu/Back-End/venv/bin/uwsgi -i /home/ubuntu/Back-End/.config/uwsgi/project.ini
 
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all
 
[Install]
WantedBy=multi-user.target